# systemd/dbd25-ntfy-test.service
[Unit]
Description=Send ntfy survey reminder %i

[Service]
Type=oneshot

# Wichtig: passe Pfade an deinen Clone-Ordner an oder benutze WorkingDirectory
WorkingDirectory=%h/ntfy-reminder

# env-file außerhalb von Git (wie oben)
EnvironmentFile=%h/ntfy-reminder/config/ntfy.env

# Aufruf: Nachricht %i senden
ExecStart=/usr/bin/env PYTHONPATH=src python3 -m 00_run send %i --env-file config/ntfy.env
